# dploi_fabric/templates/supervisor/supervisor.conf
[program:{{ process_name }}]
command={{ process_cmd }}
directory={{ deployment.path }}../
user={{ deployment.user }}
autostart=True
autorestart=True
redirect_stderr=True
environment=HOME="{{ deployment.home }}";PYTHONPATH="{{ deployment.path }}",

#################################################
