# dploi_fabric/templates/supervisor/supervisor.conf
[program:{{ process_name }}]
command={{ process_cmd }}
directory={{ deployment.path }}../
user={{ deployment.user }}
autostart=True
autorestart=True
redirect_stderr=True
environment=HOME="{{ deployment.home }}";PYTHONPATH="{{ deployment.path }}",

stdout_logfile={{ deployment.home }}/log/{{ process_name }}.stdout.log
stdout_logfile_backups=2
stderr_logfile={{ deployment.home }}/log/{{ process_name }}.stderr.log
#################################################
